# team-fullmetal/dockerComp/src/client/conf/app_supervisor.conf
[program:flaskapp]
command = uwsgi --ini /opt/dockerComp/conf/app_uwsgi.ini

[program:nginx]
command = service nginx restart
